# systemd/mnt-nas.automount
[Unit]
Description=Automount NAS

[Automount]
Where=/mnt/nas

[Install]
WantedBy=multi-user.target
