# systemd/squeeze-software.automount
[Unit]
Description=Automount Squeeze software share

[Automount]
Where=/squeeze/software

[Install]
WantedBy=multi-user.target
