# systemd/squeeze-Files.automount
[Unit]
Description=Automount Squeeze Files share

[Automount]
Where=/squeeze/Files

[Install]
WantedBy=multi-user.target
