# systemd/squeeze-Data.automount
[Unit]
Description=Automount Squeeze Data share

[Automount]
Where=/squeeze/Data

[Install]
WantedBy=multi-user.target
